# s32v234_sdk/os/linux/aarch64/vlab/etc/inittab
::sysinit:/etc/init.d/rcS

#ttyAMA0::respawn:/sbin/getty -L ttyAMA0 57600 vt100
ttyAMA0::wait:/bin/sh 

::restart:/sbin/init 
::shutdown:/bin/umount -a -r

